# ryu/ryu.service.default
[Unit]
Description=Ryu Controller
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/ryu-manager  --log-file=/var/log/ryu.log --app-lists APPS

[Install]
WantedBy=multi-user.target



