# ryu/ryu.service.default
[Unit]
Description=Ryu Controller
After=syslog.target network.target

[Service]
Type=simple
ExecStart=RYU --ofp-listen-host IP --ofp-tcp-listen-port PORT --log-file=/var/log/ryu_IP_PORT.log --observe-links --app-lists APPS EXTRAS 

[Install]
WantedBy=multi-user.target



